# net/qiperfd/files/qiperfd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2025 jimmy

START=50
STOP=15

QIPERFD=qiperfd
QIPERFD_BIN="/usr/sbin/$QIPERFD"

start() {
        echo start $QIPERFD
        # commands to launch application
        $QIPERFD_BIN &
}

stop() {
        echo stop $QIPERFD
        # commands to kill application
        killall $QIPERFD
}

